# dev-util/devpi-server/files/devpi-server.initd
#!/sbin/openrc-run

command=/usr/bin/devpi-server
command_user=devpi
command_background=1
pidfile=/run/devpi-server.pid

depend() {
	need net
}

